# systemd/postgresql@17-main.service.d/ssd-mount.conf
[Unit]
RequiresMountsFor=/srv/data/postgresql/17/main
