# systemd/airmon-agent.service
[Unit]
Description=airmon Pi sensor agent
After=network-online.target airmon-server.service
Wants=network-online.target airmon-server.service
RequiresMountsFor=/srv/data

[Service]
Type=simple
User=umut
Group=umut
SupplementaryGroups=dialout i2c gpio
WorkingDirectory=/home/umut/airmon
Environment=AIRMON_SERVER_URL=http://127.0.0.1:8000/api
Environment=AIRMON_BUFFER_DB=/srv/data/airmon/buffer.db
Environment=AIRMON_SAMPLE_INTERVAL_S=5
Environment=AIRMON_MHZ19B_DEVICE=/dev/ttyAMA0
Environment=AIRMON_SPS30_DEVICE=/dev/ttyAMA3
ExecStart=/home/umut/airmon/.venv/bin/python /home/umut/airmon/agent.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
