# systemd/airmon-server.service
[Unit]
Description=airmon FastAPI server
After=network-online.target
Wants=network-online.target
RequiresMountsFor=/srv/data

[Service]
Type=simple
User=umut
Group=umut
WorkingDirectory=/home/umut/airmon-server
Environment=AIRMON_DB_PATH=/srv/data/airmon/server.db
Environment=AIRMON_STATIC_DIR=/home/umut/airmon-web
ExecStart=/home/umut/airmon-server/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
