# systemd/homelab-backup.timer
[Unit]
Description=Nightly backups at 03:15 local

[Timer]
OnCalendar=*-*-* 03:15:00
Persistent=true
RandomizedDelaySec=15min

[Install]
WantedBy=timers.target
