# systemd/homelab-ddns.service
[Unit]
Description=Update Cloudflare A records to current public IP
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
EnvironmentFile=/etc/homelab/cloudflare.env
ExecStart=/opt/homelab/scripts/update-ddns.sh
