# systemd/homelab-backup.service
[Unit]
Description=Nightly backups: pg_dump ColdTrace + snapshot airmon SQLite, mirror offsite
After=network-online.target postgresql@17-main.service
Wants=network-online.target
RequiresMountsFor=/srv/data
ConditionPathIsMountPoint=/srv/data

[Service]
Type=oneshot
User=root
EnvironmentFile=/etc/homelab/backup.env
ExecStart=/opt/homelab/scripts/backup.sh
